# docker/supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=warn

[program:dockerd]
command=dockerd --host=unix:///var/run/docker.sock --storage-driver=overlay2
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dockerd.log
stderr_logfile=/var/log/supervisor/dockerd-error.log
priority=10
startsecs=3

[program:gateway]
command=/bin/sh -lc 'exec hermes gateway run --replace'
directory=/workspace
autostart=true
autorestart=true
stdout_logfile=/data/hermes/logs/gateway-supervisor.log
stderr_logfile=/data/hermes/logs/gateway-supervisor-error.log
priority=20
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
